# server/brpaste.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

respawn_delay=30

command_user="nobody:nobody"
command=/usr/bin/brpaste
command_args="\
    ${BRPASTE_REDIS:+-redis $BRPASTE_REDIS} \
    ${BRPASTE_ADDR:+-bind $BRPASTE_ADDR}"

depend() {
    need redis
}
